# scripts/com.stocktrading.dataupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd 排程：平日 14:30 盤後自動增量更新市場數據。
     - 股價走 shioaji（快，數分鐘完成）
     - 法人/融資券/月營收/除權息走 FinMind，額度用罄自動等下個整點續跑
       （--auto-wait，過夜也會自己補完，最多等 30 個小時週期）
     安裝：
       cp scripts/com.stocktrading.dataupdate.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.stocktrading.dataupdate.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stocktrading.dataupdate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/omar/myprojects/StockTradingSystem/.venv/bin/python</string>
        <string>-m</string>
        <string>scripts.backfill</string>
        <string>--auto-wait</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/omar/myprojects/StockTradingSystem</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/omar/myprojects/StockTradingSystem/logs/dataupdate_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/omar/myprojects/StockTradingSystem/logs/dataupdate_launchd.log</string>
</dict>
</plist>
